# package/src/deb-scripts/postinst
cp /usr/local/share/avm-vdriver/avm-vdriver.service /lib/systemd/system/avm-vdriver.service
chmod 644 /lib/systemd/system/avm-vdriver.service
systemctl daemon-reload
systemctl enable avm-vdriver.service
systemctl start avm-vdriver.service
systemctl status avm-vdriver.service --no-pager || true
