# package/src/deb-scripts/prerm
systemctl disable avm-vdriver.service 2>/dev/null
systemctl stop avm-vdriver.service 2>/dev/null
